# zdyndns_client.conf
[program:consul]
command=/home/ubuntu/zdyndns_client.sh
autostart=true
autorestart=true

